# etc/libup-run.timer
[Unit]
Description=Run the script

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target

